# contrib/udev/40-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/etc/pm/sleep.d/40autorandr thaw udev"
