# contrib/udev/40-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/bin/autorandr --batch -c --default default"
